# bunker-addon/rootfs/etc/services.d/bunker/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Bunker
# Starts the Bunker status page
# ==============================================================================

# Set the current working dir to /opt
cd /opt/bunker || bashio::exit.nok "Could not change directory to Bunker"

bashio::log.info "Starting Bunker..."

# And finally start the Bunker
exec node app.js