# happy-addon/rootfs/etc/services.d/happy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Home Assistant Add-ons: Happy
# Starts the Happy status page
# ==============================================================================

# Set the current working dir to /opt
cd /opt/happy || bashio::exit.nok "Could not change directory to Happy"

bashio::log.info "Starting Happy..."

# And finally start the Happy
exec node app.js