# filiikot-addon/rootfs/etc/services.d/filiikot/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: FiliiKot
# Starts the FiliiKot status page
# ==============================================================================

# Set the current working dir to /opt
cd /opt/filiikot || bashio::exit.nok "Could not change directory to FiliiKot"

bashio::log.info "Starting FiliiKot..."

# And finally start the FiliiBot
exec node app.js