# filiibot-addon/rootfs/etc/services.d/filiibot/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: FiliiBot
# Starts the FiliiBot
# ==============================================================================

cd /opt/filiibot || bashio::exit.nok "Could not change directory to FiliiBot"

bashio::log.info "Starting FiliiBot..."

# Starts the FiliiBot
exec node bot.js