# filiicodex-addon/rootfs/etc/services.d/filiibot/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: FiliiCodex
# Starts the FiliiCodex
# ==============================================================================

# Set the current working dir to /opt
cd /opt/filiicodex || bashio::exit.nok "Could not change directory to FiliiCodex"

bashio::log.info "Starting FiliiCodex..."

# And finally start the FiliiCodex
exec node index.js